# packaging/systemd/icahazip.service
[Unit]
Description=Icahanzip golang implementation
Documentation=https://github.com/charliemaiors/icahazip
After=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/icanhazip --config /etc/icahazip/config.yaml
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=icahazip

# Resource limits
LimitNOFILE=8192
MemoryMax=256M

[Install]
WantedBy=multi-user.target